# infra/edge/systemd/safegai-edge.service
[Unit]
Description=SafeGAI Edge Gateway
After=network-online.target mediamtx.service
Wants=network-online.target

[Service]
Type=notify
User=safegai
Group=safegai
ExecStart=/opt/safegai/current/bin/safegai-edge --config /etc/safegai/config.yaml
Restart=always
RestartSec=3
WatchdogSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/safegai /var/log/safegai
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
